# services/dashcam-monitor/run
#!/bin/sh
exec /usr/bin/setuidgid pi /home/pi/teslacam/src/dashcam-monitor.js
