# services/dashcam-monitor/run
#!/bin/sh
exec /home/pi/teslacam/src/dashcam-monitor.js
