# services/rsync-upload/run
#!/bin/sh
exec /home/pi/teslacam/src/uploaders/rsync.js
