# services/rsync-upload/run
#!/bin/sh
exec setuidgid pi /home/pi/teslacam/src/uploaders/rsync.js
