# services/rsync-upload/run
#!/bin/sh
exec /usr/bin/setuidgid pi /home/pi/teslacam/src/uploaders/rsync.js
